# ansible/roles/certpackage/files/certpackage.service
[Unit]
Description=Package Lets-Encrypt Certificates for sync.

[Service]
Type=oneshot
ExecStart=/usr/local/bin/certpackage.sh
PrivateTmp=true
