# cookbooks/mounts/templates/unit.mount.erb
[Unit]
Description=Mount <%= @mount_name %> Share at boot
After=network-online.target
Wants=network-online.target

[Mount]
What=//<%= node['configs']['mount']['target_server'] %>/<%= @mount_name %>
Where=/mnt/<%= @mount_name %>
Options=credentials=<%= node['configs']['mount']['samba_credentials_path'] %>/<%= @mount_name %>,iocharset=utf8,rw,x-systemd.automount
Type=cifs
TimeoutSec=30

[Install]
WantedBy=multi-user.target
