# cookbooks/notifications/files/failure-notification@.service
[Unit]
Description=Failure emergency notification service.

[Service]
Type=oneshot
ExecStart=/usr/local/bin/failure_notification.sh %i
