# certpackage/systemd/certpackage.service
[Unit]
Description=Certpackage

[Service]
Type=oneshot
ExecStart=/usr/local/bin/certpackage.sh
PrivateTmp=true

